# net-p2p/monero/files/init
#!/sbin/openrc-run

depend() {
    need net
    use dns
}

start() {
    ebegin "Starting monero daemon"
    checkpath -d -m 0755 -o monero:monero /var/run/monero
    checkpath -d -m 0755 -o monero:monero ${DATA_DIR}
    start-stop-daemon --start --quiet --background --exec monerod \
    --user ${MONEROD_USER} --group ${MONEROD_GROUP} --pidfile ${PID_FILE} \
    -- --non-interactive --detach --pidfile ${PID_FILE} --data-dir ${DATA_DIR}
    eend $?
}

stop() {
    ebegin "Stopping monero daemon"
    start-stop-daemon --stop --exec monerod --quiet \
    --user ${MONEROD_USER} --group ${MONEROD_GROUP} --pidfile ${PID_FILE}
    eend $?
}
