# net-p2p/scavenger/files/init
#!/sbin/openrc-run

start(){
    ebegin "Starting Burst miner"
    start-stop-daemon --background --start --exec /usr/bin/scavenger \
    --make-pidfile --pidfile /var/run/scavenger.pid \
    -u ${MINER_USER} -- -c /etc/scavenger/config.yaml

    eend $?
}

stop(){
    ebegin "Stopping Burst miner"
    start-stop-daemon --stop --exec /usr/bin/scavenger \
    --pidfile /var/run/scavenger.pid
    eend $?
}