# www-apps/invidious/files/init
#!/sbin/openrc-run

depend() {
    need net
    after postgresql-14
}

start() {
    ebegin "Starting Invidious"
    checkpath -d -m 0755 -o invidious:invidious /var/run/monero
    checkpath -d -m 0755 -o invidious:invidious ${DATA_DIR}
    start-stop-daemon --start --quiet --background --exec invidious \
                      --user ${INVID_USER} --group ${INVID_GROUP} \
                      --make-pidfile --pidfile ${PID_FILE} \
                      --chdir ${DATA_DIR}
    eend $?
}

stop() {
    ebegin "Stopping Invidious"
    start-stop-daemon --stop --exec monerod --quiet \
                      --user ${INVID_USER} --group ${INVID_GROUP} \
                      --pidfile ${PID_FILE}
    eend $?
}
