# net-p2p/xmr-stak-amd/files/init
#!/sbin/openrc-run

start(){
    ebegin "Starting Monero CPU miner"
    start-stop-daemon --background --start --exec /usr/bin/xmr-stak-amd \
    --make-pidfile --pidfile /var/run/xmr-stak-amd.pid \
    -u ${MINER_USER} \
    -- /etc/xmr-stak/amd/config.txt
    eend $?
}

stop(){
    ebegin "Stopping Monero CPU miner"
    start-stop-daemon --stop --exec /usr/bin/xmr-stak-amd \
    --pidfile /var/run/xmr-stak-amd.pid
    eend $?
}