# net-p2p/burstcoin/files/init
#!/sbin/openrc-run

depend() {
    need net mysql
}

start() {
    ebegin "Starting burst wallet daemon"
    checkpath -d -m 0755 -o ${BURST_USER}:${BURST_GROUP} /var/run/burst
    start-stop-daemon --start --quiet --background \
    --chdir /opt/burst/ \
    --user ${BURST_USER} --group ${BURST_GROUP} \
    --pidfile ${PID_FILE} --make-pidfile \
    --exec java -- -cp burst.jar:conf brs.Burst
    eend $?
}

stop() {
    ebegin "Stopping burst wallet daemon"
    start-stop-daemon --stop --chdir /opt/burst --exec /opt/burst/burst.sh --quiet \
    --user ${BURST_USER} --group ${BURST_GROUP} --pidfile ${PID_FILE}
    eend $?
}
